# systemd/auto-update.service
[Unit]
Description=Daily system update and upgrade
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
ExecStart=/opt/vps-ops-toolkit/scripts/auto-update.sh
